# crates/app-cli/src/daemon/templates/systemd.service
[Unit]
Description=repo-link reconciliation daemon
After=network-online.target

[Service]
ExecStart={{BINARY_PATH}} --log-format=json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
